# deploy/k8s/supervisord.conf
; supervisord configuration for hKask pod.
; Manages three processes: Conduit (Matrix), Litestream (WAL backup), kask (main app).

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:conduit]
command=/usr/local/bin/conduit
environment=CONDUIT_CONFIG="/etc/conduit/conduit.toml"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:litestream]
command=/usr/local/bin/litestream replicate -config /etc/litestream.yml
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:kask]
command=/usr/local/bin/kask serve
environment=POD_ID="%(ENV_POD_ID)s",HKASK_DATA_DIR="/data",HKASK_BASE_URL="%(ENV_HKASK_BASE_URL)s",HKASK_DOMAIN="%(ENV_HKASK_DOMAIN)s",HKASK_MATRIX_URL="http://localhost:8008",HKASK_OAUTH_GITHUB_CLIENT_ID="%(ENV_HKASK_OAUTH_GITHUB_CLIENT_ID)s",HKASK_OAUTH_GITHUB_CLIENT_SECRET="%(ENV_HKASK_OAUTH_GITHUB_CLIENT_SECRET)s"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
